# jenkins-slaves/ansible/data/service/verdaccio.service
[Unit]
Description=Verdaccio Service.

[Service]
Type=simple

ExecStart=/usr/local/bin/verdaccio.service.sh

[Install]
WantedBy=multi-user.target
