# jenkins-agents/ansible/data/service/smee.service
[Unit]
Description=Smee Client Service.

[Service]
Type=simple

ExecStart=/usr/local/bin/smee.service.sh

[Install]
WantedBy=multi-user.target
